# supervisor/sshd.conf
[program:sshd]
command=/usr/sbin/sshd -D
priority=0
